# debian/postinst
#!/bin/sh -e

if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-remove" ] ; then
    #add brightness-init script to restore brightness after reboot
    if [ -e /etc/init.d/brightness-init ]; then
    update-rc.d brightness-init remove
    rm /etc/init.d/brightness-init
    fi
fi

#DEBHELPER#
